# templates/default/services/gem-mirror.conf.erb
#!upstart
description "Update gem mirror"

start on started gem-mirror-shim
stop on stopping gem-mirror-shim

task
respawn

setuid <%= node.gem_mirror.user %>

env USER=<%= node.gem_mirror.user %>
env HOME=/home/<%= node.gem_mirror.user %>
chdir /home/<%= node.gem_mirror.user %>

pre-start script
  . /etc/environment
  bundle install
end script

exec bundle exec gem mirror -v
